# supervisord.conf
[supervisord]
nodaemon=true
user=frappe
pidfile=/home/frappe/supervisord.pid
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/frappe/supervisor.sock

[unix_http_server]
file=/home/frappe/supervisor.sock
chmod=0700

# ─── 1. Gunicorn Backend (listens locally on port 8000, proxied by Nginx) ───
[program:gunicorn]
command=/home/frappe/frappe-bench/env/bin/gunicorn -b 127.0.0.1:8000 frappe.app:application --preload --workers 2 --timeout 120
directory=/home/frappe/frappe-bench/sites
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── 2. Node Socket.io Server ───
[program:socketio]
command=node /home/frappe/frappe-bench/apps/frappe/socketio.js
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── 3. Frappe Scheduler ───
[program:scheduler]
command=/usr/local/bin/bench schedule
directory=/home/frappe/frappe-bench
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── 4. Short Worker ───
[program:worker_short]
command=/usr/local/bin/bench worker --queue short,default
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── 5. Long Worker ───
[program:worker_long]
command=/usr/local/bin/bench worker --queue long,default,short
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ─── 6. Nginx Frontend (listens on 8080 publicly, proxies to Gunicorn on 8000) ───
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
directory=/home/frappe/frappe-bench
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
