# supervisord.conf
[supervisord]
nodaemon=true
user=frappe
pidfile=/home/frappe/supervisord.pid
logfile=/home/frappe/frappe-bench/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/frappe/supervisor.sock

[unix_http_server]
file=/home/frappe/supervisor.sock
chmod=0700

# ─── 1. Gunicorn Backend ───
[program:gunicorn]
command=/home/frappe/frappe-bench/env/bin/gunicorn -b 127.0.0.1:8000 frappe.app:application --preload --workers 3 --timeout 120
directory=/home/frappe/frappe-bench/sites
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/gunicorn.log
stderr_logfile=/home/frappe/frappe-bench/logs/gunicorn.err.log

# ─── 2. Node Socket.io Server ───
[program:socketio]
command=node /home/frappe/frappe-bench/apps/frappe/socketio.js
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/socketio.log
stderr_logfile=/home/frappe/frappe-bench/logs/socketio.err.log

# ─── 3. Frappe Scheduler ───
[program:scheduler]
command=/home/frappe/frappe-bench/env/bin/bench schedule
directory=/home/frappe/frappe-bench
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/scheduler.log
stderr_logfile=/home/frappe/frappe-bench/logs/scheduler.err.log

# ─── 4. Short Worker ───
[program:worker_short]
command=/home/frappe/frappe-bench/env/bin/bench worker --queue short,default
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/worker_short.log
stderr_logfile=/home/frappe/frappe-bench/logs/worker_short.err.log

# ─── 5. Long Worker ───
[program:worker_long]
command=/home/frappe/frappe-bench/env/bin/bench worker --queue long,default,short
directory=/home/frappe/frappe-bench
environment=FRAPPE_REDIS_CACHE="redis://%(ENV_REDIS_HOST)s:6379/0",FRAPPE_REDIS_QUEUE="redis://%(ENV_REDIS_HOST)s:6379/1"
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/worker_long.log
stderr_logfile=/home/frappe/frappe-bench/logs/worker_long.err.log

# ─── 6. Nginx Frontend ───
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
directory=/home/frappe/frappe-bench
autostart=true
autorestart=true
stdout_logfile=/home/frappe/frappe-bench/logs/nginx.log
stderr_logfile=/home/frappe/frappe-bench/logs/nginx.err.log
